# contrib/systemd/ruvname-default-config.service
[Unit]
Description=RUVNAME default config generator
ConditionPathExists=|!/etc/ruvname.toml
ConditionFileNotEmpty=|!/etc/ruvname.toml
Wants=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
User=ruvname
Group=ruvname
StandardOutput=file:/etc/ruvname.toml
ExecStart=/usr/bin/ruvname -g
